# packaging/linux/templates/usr/lib/systemd/system/audex.service
[Unit]
Description=Audex Medical Transcription Service
Documentation=https://github.com/6ixGODD/audex
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
SupplementaryGroups=audio
WorkingDirectory=/usr/lib/audex

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="ALSA_PCM_CARD=0"
Environment="ALSA_PCM_DEVICE=0"
Environment="PULSE_RUNTIME_PATH=/tmp/audex-audio"
Environment="PULSE_SERVER=unix:/tmp/audex-pulse.sock"
Environment="PA_ALSA_PLUGHW=1"

# Ensure directories exist
ExecStartPre=/bin/mkdir -p /var/log/audex
ExecStartPre=/bin/mkdir -p /var/lib/audex/store
ExecStartPre=/bin/mkdir -p /etc/audex
ExecStartPre=/bin/mkdir -p /tmp/audex-audio
ExecStartPre=/bin/chmod 777 /tmp/audex-audio

# Start Audex
ExecStart=/usr/bin/audex

# Restart policy
Restart=on-failure
RestartSec=10s

# Security
NoNewPrivileges=false
PrivateTmp=false

# Audio access
SupplementaryGroups=audio

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audex

[Unit]
StartLimitIntervalSec=300
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
