# packaging/linux/templates/usr/lib/systemd/system/audex.service
[Unit]
Description=Audex Medical Transcription Service for %I
Documentation=https://github.com/6ixGODD/audex
After=network-online.target sound.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/audex

# Wait for user session to be ready
ExecStartPre=/bin/sleep 5

# Start Audex as user
ExecStart=/usr/bin/audex

# Restart policy
Restart=on-failure
RestartSec=30s

# Environment
Environment="XDG_RUNTIME_DIR=/run/user/%U"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"

# Audio access (user should be in audio group)
SupplementaryGroups=audio

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audex-%i

[Unit]
StartLimitIntervalSec=300
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
