# packaging/linux/templates/DEBIAN/prerm
#!/bin/bash
# ============================================================================
# Pre-removal script for Audex DEB package
# ============================================================================

set -e

# ============================================================================
# Colors
# ============================================================================

BLUE='\033[0;34m'
NC='\033[0m'

# ============================================================================
# Handle Removal
# ============================================================================

case "$1" in
    remove|upgrade|deconfigure)
        echo -e "${BLUE}ðŸ—‘ï¸  Removing Audex...${NC}"

        # Stop systemd service if running
        if systemctl is-active --quiet audex.service 2>/dev/null; then
            echo "   Stopping audex service..."
            systemctl stop audex.service 2>/dev/null || true
        fi

        # Disable systemd service if enabled
        if systemctl is-enabled --quiet audex.service 2>/dev/null; then
            echo "   Disabling audex service..."
            systemctl disable audex.service 2>/dev/null || true
        fi

        echo "âœ“ Audex service stopped"
        ;;
esac

exit 0
