# packaging/linux/templates/DEBIAN/prerm
#!/bin/bash
# ============================================================================
# Pre-removal script for Audex DEB package
# ============================================================================

set -e

# ============================================================================
# Colors
# ============================================================================

BLUE='\033[0;34m'
YELLOW='\033[1;33m'
NC='\033[0m'

# ============================================================================
# Handle Removal
# ============================================================================

case "$1" in
    remove|upgrade|deconfigure)
        echo -e "${BLUE}ðŸ—‘ï¸  Removing Audex... ${NC}"

        # Stop and disable system-level services for all users
        for user_service in /etc/systemd/system/audex@*.service; do
            if [ -f "$user_service" ]; then
                service_name=$(basename "$user_service")
                echo "   Stopping $service_name..."
                systemctl stop "$service_name" 2>/dev/null || true
                systemctl disable "$service_name" 2>/dev/null || true
            fi
        done

        # Find all users with audex service enabled
        while IFS=: read -r user_name _ user_id _ _ user_home _; do
            if [ "$user_id" -ge 1000 ] && [ "$user_id" -lt 65534 ]; then
                if [ -d "$user_home/.config/systemd/user" ]; then
                    # Stop user-level service if exists
                    sudo -u "$user_name" XDG_RUNTIME_DIR="/run/user/$user_id" \
                        systemctl --user stop audex.service 2>/dev/null || true
                    sudo -u "$user_name" XDG_RUNTIME_DIR="/run/user/$user_id" \
                        systemctl --user disable audex.service 2>/dev/null || true
                fi
            fi
        done < /etc/passwd

        echo "âœ“ Audex services stopped"
        ;;
esac

exit 0
